# desktop/rootfs/etc/supervisor/conf.d/pulseaudio.conf
[program:pulseaudio]
directory=/home/desktop
command=/usr/bin/pulseaudio
user=desktop
environment=DISPLAY=":1",HOME="/home/desktop",USER="desktop"
